# config/supervisord/alarm.conf
[program:alarm]
directory=/opt/alarm/backend
user=root
command=/opt/alarm/env.sh gunicorn -c gunicorn_config.py topotal.wsgi
stdout_logfile=/var/log/supervisord.log
stderr_logfile=/var/log/supervisord.log
#stdout_logfile=/dev/null
#stderr_logfile=/dev/null
autostart=true
redirect_stderr=true
